# roles/open-webui/files/com.sloria.open-webui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>com.sloria.open-webui</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>tool</string>
        <string>run</string>
        <string>--python</string>
        <string>3.11</string>
        <string>open-webui@latest</string>
        <string>serve</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DATA_DIR</key>
        <string>/opt/homebrew/var/open-webui</string>
        <key>WEBUI_NAME</key>
        <string>RC LLM</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/opt/homebrew/var/open-webui/open-webui.log</string>
    <key>StandardOutPath</key>
    <string>/opt/homebrew/var/open-webui/open-webui.log</string>
    <key>WorkingDirectory</key>
    <string>/opt/homebrew/var/open-webui</string>
</dict>
</plist>
